# aware-filter.service
[Unit]
Description=AWARE Webservice Receiver
After=network.target

[Service]
Type=simple
User=rantahj1
WorkingDirectory=%h/aware_filter
EnvironmentFile=-%h/aware_filter/.env.service
ExecStart=%h/aware_filter/venv/bin/python -m aware_filter
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target